# support/gps/udev.d/gps.rules
# Adafruit Ultimate GPS HAT
KERNEL=="ttyS0", RUN+="/bin/stty -F /dev/%k raw 9600 cs8 clocal -cstopb"
KERNEL=="pps0", OWNER="gpsd"
